# config/gateway.conf
# Upstart Job to execute the Spring Boot application on startup
description "gateway"
start on startup
task
exec /vagrant/gateway/config/spring-start.sh
